# systemd/ledmatrix.service
[Unit]
Description=LED Matrix Display Service
After=network.target

[Service]
Type=simple
User=ledpi
Group=ledpi
WorkingDirectory=__PROJECT_ROOT_DIR__
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=__VENV_PYTHON__ __PROJECT_ROOT_DIR__/run.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ledmatrix
# Support for on-demand plugin filtering via environment variable
# The environment variable LEDMATRIX_ON_DEMAND_PLUGIN can be set via:
# sudo systemctl set-environment LEDMATRIX_ON_DEMAND_PLUGIN=<plugin_id>
# Or by using an EnvironmentFile (see below)
# EnvironmentFile=__PROJECT_ROOT_DIR__/config/on_demand_env.conf

[Install]
WantedBy=multi-user.target 